# net/tunnel/cmd/x.tunnel.server.supervisor.conf
[program:x-tunnel-server]
command=/usr/local/bin/x.tunnel.server -d
directory=/tmp
user=root
autostart=true
autorestart=true