# net/tunnel/cmd/x.tunnel.client.supervisor.conf
[program:x-tunnel-client]
command=/usr/local/bin/x.tunnel.client
directory=/tmp
user=root
autostart=true
autorestart=true