# net/tunnel/cmd/client/supervisor.conf
[program:gox.tunnel.client]
command=/usr/local/bin/gox.tunnel.client
directory=/tmp
user=root
autostart=true
autorestart=true