# net/tunnel/cmd/client/supervisor.conf
[program:tunnel.client]
command=/usr/local/bin/tunnel.client
directory=/tmp
user=root
autostart=true
autorestart=true