# net/tunnel/cmd/server/supervisor.conf
[program:tunnel.server]
command=/usr/local/bin/tunnel.server
directory=/tmp
user=root
autostart=true
autorestart=true