# net/tunnel/cmd/server/supervisor.conf
[program:gox.tunnel.server]
command=/usr/local/bin/gox.tunnel.server
directory=/tmp
user=root
autostart=true
autorestart=true